# modules/gitorious/templates/etc/init/gitorious-unicorn.conf.erb
description "Run unicorn process"
author "Paweł Pierzchała <pawel@gitorious.org>"

start on started rc RUNLEVEL=[35]
stop on runlevel [06]

env PATH=/bin:/usr/bin:/usr/local/bin

exec /usr/local/bin/chruby-exec ruby-<%= @ruby_version %> -- <%= @app_root %>/bin/unicorn
respawn
