# roles/gitorious/templates/upstart/gitorious-sphinx.conf
description "Gitorious search engine (sphinx)"
author "Marcin Kulik"

start on runlevel [2345]
stop on runlevel [!2345]

exec su - git -c "exec /var/www/gitorious/app/bin/rake ts:start_in_foreground >>/var/log/gitorious/sphinx.log 2>&1"

respawn
