# modules/resque/templates/etc/init/resque-worker.conf.erb
description "Run a Resque worker on all queues"
author "Marius Mårnes Mathiesen <marius@gitorious.com>"

start on started rc RUNLEVEL=[35]
stop on runlevel [06]

env PATH=/bin:/usr/bin:/usr/local/bin
env QUEUE=*
env PIDFILE=<%= resque_gitorious_root %>/tmp/pids/resque-worker1.pid

exec /usr/local/bin/chruby-exec ruby-1.9.3-p448 -- <%= resque_gitorious_root %>/bin/rake resque:work
respawn