# upstart/gitorious-redis.conf
description "gitorious-redis"
author "Marcin Kulik"

start on filesystem and started docker
stop on runlevel [!2345]
respawn

script
  /usr/bin/docker start -a redis
end script
