# upstart/gitorious-memcached.conf
description "gitorious-memcached"
author "Marcin Kulik"

start on filesystem and started docker
stop on runlevel [!2345]
respawn

script
  /usr/bin/docker start -a memcached
end script
