# upstart/gitorious-git-daemon.conf
description "gitorious-git-daemon"
author "Marcin Kulik"

start on filesystem and started docker
stop on runlevel [!2345]
respawn

script
  /usr/bin/docker start -a git-daemon
end script
